# wallpapers/.config/systemd/user/wallpaper-cycle.service
[Unit]
Description=Wallpaper cycle (theme+time based)

[Service]
Type=oneshot
ExecStart=%h/.wallpapers/bin/wallpaper_cycle.py
